# deployment/chef/src/main/resources/cookbooks/cadmium/templates/default/cadmium-jetty.conf.erb
start on runlevel [2345]
stop on runlevel [016]

respawn
chdir <%= @cwd %>
setuid <%= @user %>
setgid <%= @group %>
exec java -server -Djava.awt.headless=true -Djava.net.preferIPv4Stack=true \
  -Djboss.server.log.dir=./log \
  -Dcom.meltmedia.cadmium.environment=<%= @environment %> \
  -Dcom.meltmedia.cadmium.contentRoot=<%= @content_root %> \
  -Dport=<%= @port %> \
  -jar ./bin/<%= @war_name %>